# cluster/mesos-master.conf
[program:mesos-master]
command=/usr/bin/mesos-init-wrapper master
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
