# cluster/mesos-slave.conf
[program:mesos-slave]
command=/usr/bin/mesos-init-wrapper slave
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
